# ansible/roles/configure_k8s_master/templates/etc/profile.d/domain_env.sh
#!/bin/sh
export PREFIX="{{ k8s_prefix }}"
export DOMAIN="{{ k8s_domain }}"
export RELEASENAME="{{ releasename }}"
export NAMESPACE="{{ namespace }}"
export KEYFILE="{{ keyfile }}"
export CERTFILE="{{ certfile }}"
